# systemd/trading-bot.service
[Unit]
Description=Hyperliquid Trading Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=trading
Group=trading
WorkingDirectory=/home/trading/tradingbot
EnvironmentFile=/home/trading/tradingbot/.env
ExecStart=/home/trading/tradingbot/.venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trading-bot

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/home/trading/tradingbot

# Resource limits
LimitNOFILE=65536
MemoryMax=2G

[Install]
WantedBy=multi-user.target
